# cookbooks/seedbox/files/default/homepage.conf
[program:homepage]
command=python /opt/homepage/homepage.py
autostart=true
autorestart=true
